# systemd/hermes-agent.service
[Unit]
Description=Hermes Agent Docker Compose stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/hermes-agent
ExecStart=/usr/bin/docker compose -f /opt/hermes-agent/docker-compose.yml up -d
ExecStop=/usr/bin/docker compose -f /opt/hermes-agent/docker-compose.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
